# roles/minikube/templates/kubectl_proxy.service.j2
[Unit]
Description=kubectl proxy

[Service]
User=marco
ExecStart=/etc/minikube/run_kubectl_proxy.sh

[Install]
WantedBy=multi-user.target
